# crontab/crontab
05 01 * * 0 otrade --strategy Simple >> /var/log/cron.log 2>&1